# network-monitoring/nav/manager/container/target/etc/supervisor/conf.d/rsyslog.conf
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autorestart=true
stdout_logfile = /var/log/rsyslog/rsyslog-supervisor.log
redirect_stderr = true
