# networks/monitoring/nav/manager/container/target/etc/nav/cron.d/netbiostracker
## info: Regularly fetch netbiosnames from active computers
*/15 * * * * netbiostracker.py >> /var/log/cron/netbiostracker.log 2>&1
