# .playground/nav/container/etc/supervisor/conf.d/cron.conf
[program:cron]
command = /usr/sbin/cron -f -L 15
autorestart = true
stdout_logfile = /var/log/cron/cron-stdout.log
stderr_logfile = /var/log/cron/cron-stderr.log
