# network-monitoring/nav/container/target/etc/nav/cron.d/maintengine
*/5 * * * * /usr/local/bin/maintengine.py
