# network-monitoring/nav/container/target/etc/nav/cron.d/thresholdmon
## info: Monitors your Graphite metrics for exceeded thresholds
*/5 * * * * thresholdmon >> /var/log/cron/thresholdmon.log 2>&1
