# nav/container/etc/nav/cron.d/mactrace
11,26,41,56 * * * * /usr/local/bin/macwatch.py
