# nav/container/target/etc/nav/cron.d/activeip
*/30 * * * * /usr/local/bin/collect_active_ip.py
