# nav/container/target/etc/nav/cron.d/topology
35 * * * * /usr/local/bin/navtopology --l2 --vlan
