# .playground/nav/container/etc/nav/cron.d/dbclean
*/5 * * * * /usr/local/bin/navclean -f -q --netbox --websessions
