# .playground/nav/container/etc/nav/cron.d/emailreports
0 5 1 * * /usr/local/bin/emailreports.py month
0 5 * * mon /usr/local/bin/emailreports.py week
0 5 * * * /usr/local/bin/emailreports.py day
