# .playground/nav/container/etc/nav/cron.d/logengine
* * * * * /usr/local/bin/logengine.py -q
3 3 * * * /usr/local/bin/logengine.py -d
3 * * * * /usr/local/bin/ipdevpolld --clean
