# network-monitoring/nav/manager/container/target/etc/nav/cron.d/mactrace
## info: Checks NAV's cam log for watched MAC addresses
11,26,41,56 * * * * macwatch.py >> /var/log/cron/mactrace.log 2>&1
